# config/supervisord/supervisord.d/crond.ini
[program:crond]
command = /usr/sbin/crond -n -s
user = root
autostart = true
stdout_logfile = /data/log/supervisor/crond-output.log
stdout_logfile_maxbytes = 1MB
stderr_logfile = /data/log/supervisor/crond-error.log
stderr_logfile_maxbytes = 1MB