# config/supervisord/supervisord.d/nginx.ini
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
user = root
autostart = true
autorestart = true
startretries = 5
stdout_logfile = /data/log/supervisor/nginx-output.log
stdout_logfile_maxbytes = 10MB
stderr_logfile = /data/log/supervisor/nginx-error.log
stderr_logfile_maxbytes = 10MB