# config/supervisord/supervisord.d/php-fpm.ini
[program:php-fpm]
command = /usr/local/php-7.2.10/sbin/php-fpm --fpm-config /usr/local/php-7.2.10/etc/php-fpm.conf -F
user = root
autostart = true
autorestart = true
startretries = 5
stdout_logfile = /data/log/supervisor/php-fpm-output.log
stdout_logfile_maxbytes = 10MB
stderr_logfile = /data/log/supervisor/php-fpm-error.log
stderr_logfile_maxbytes = 10MB